# systemd/sparky-eth.service
[Unit]
Description=Sparky AI Research — ETH Strategies (Agent A)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=3600
StartLimitBurst=10

[Service]
Type=simple
WorkingDirectory=/home/akamath/sparky-ai
ExecStart=/home/akamath/sparky-ai/.venv/bin/python -m sparky.workflow.runner directives/20260218_eth_strategies.yaml
Restart=on-failure
RestartSec=120
Environment="PATH=/home/akamath/sparky-ai/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="VIRTUAL_ENV=/home/akamath/sparky-ai/.venv"

[Install]
WantedBy=default.target
